# tests/acceptance/files/gkeepd.service
[Unit]
Description=git-keeper server

[Service]
Type=simple
ExecStart=/usr/local/bin/gkeepd

[Install]
WantedBy=default.target
